# src/services/bot_hook.service
[Unit]
Description=Webhook to send telegram messages

[Service]
Type=simple
ExecStart=/usr/bin/simpleci update
