# src/services/simpleci.service
[Unit]
Description=Simple CI for AD CTF

[Service]
Type=simple
ExecStart=/usr/bin/simpleci update
