# src/services/simpleci.timer
[Unit]
Description=Run simpleci every 5 seconds

[Timer]
OnUnitActiveSec=5s

[Install]
WantedBy=timers.target