# src/services/simpleci.service
[Unit]
Description=Simple CI for AD CTF

[Service]
Type=simple
ExecStart=/usr/local/bin/simpleci check
